# Computer/EnvironmentSetup/Mac/搭建 Mac OS 开发环境/crontab/crontab.sh
*/1 * * * * /bin/sh  /Users/liuyu/Documents/git/document/Study/Other/Script/sh/auto_push_git_repositories.sh >> /dev/null # /Users/liuyu/sh.txt
